# dadude-agent/dadude-git-watchdog.service
[Unit]
Description=DaDude Git Update Watchdog
Documentation=https://github.com/grandir66/Dadude
After=network-online.target
Wants=network-online.target

# NON dipende da dadude-agent - deve funzionare anche se l'agent è corrotto
# After=dadude-agent.service

[Service]
Type=simple
User=root
Group=root

# Directory di lavoro
WorkingDirectory=/opt/dadude-agent

# Variabili d'ambiente
Environment="PYTHONUNBUFFERED=1"
Environment="AGENT_DIR=/opt/dadude-agent"
Environment="UPDATE_CHECK_INTERVAL=3600"
Environment="AGENT_HEALTH_TIMEOUT=120"
Environment="GIT_REMOTE=origin"
Environment="GIT_BRANCH=main"

# File .env opzionale per override
EnvironmentFile=-/opt/dadude-agent/.env.watchdog

# Comando di avvio (usa Python di sistema, non venv dell'agent)
ExecStart=/usr/bin/python3 /opt/dadude-agent/git-update-watchdog.py

# Restart policy - molto aggressiva per garantire che il watchdog sia sempre attivo
Restart=always
RestartSec=30

# Se fallisce troppe volte, attende di più prima di riprovare
StartLimitBurst=5
StartLimitIntervalSec=300

# Limiti risorse minimi
MemoryMax=128M
CPUQuota=10%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dadude-git-watchdog

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

# Permessi di scrittura necessari
ReadWritePaths=/opt/dadude-agent
ReadWritePaths=/var/log

# Timeout
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
